# etc/periodic/daily/looper
#!/bin/sh

cd /opt/app
s6-setuidgid "$PUID:$PGID" ./looper.sh
