# etc/s6/api/run
#!/bin/sh

cd /opt/app
s6-setuidgid "$PUID:$PGID" node ./dist/api/index.js
