# etc/s6/api/run
#!/bin/sh

set -e

cd /opt/app
exec s6-setuidgid "${PUID}:${PGID}" node ./dist/api/index.js
