# etc/s6/cron/run
#!/bin/sh

set -e

cat << EOF > /etc/periodic/daily/looper
#!/bin/sh

cd /opt/app
s6-setuidgid "${PUID}:${PGID}" ./looper.sh
EOF

chmod +x /etc/periodic/daily/looper

exec crond -s -n
